# files/ffmpeg/ffmpeg.service
[Unit]
Description=ffmpeg service

[Service]
ExecStart=/opt/ffmpeg.sh
ExecReload=/bin/kill -HUP $MAINPID
KillMode=control-group
KillSignal=SIGKILL
Restart=always

[Install]
WantedBy=multi-user.target
